# content/etc/cron.weekly/xbian-weekly-btrfs
#!/bin/sh

. /etc/default/xbian-snap

[ "$ENABLEDSCHEDULE" = yes ] || exit 0
[ "$KEEPWEEKS" -ne 0 ] || exit 0

btrfs-auto-snapshot -l weekly -k $KEEPWEEKS -v //
