# content/etc/cron.weekly/xbian-weekly-btrfs
#!/bin/sh

btrfs-auto-snapshot -l weekly -k 5 -v //
