# content/etc/cron.daily/xbian-daily-btrfs
#!/bin/sh

PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

btrfs-auto-snapshot -l daily -k 30 -v //
