# content/etc/cron.daily/xbian-daily-btrfs
#!/bin/sh

btrfs-auto-snapshot -l daily -k 30 -v //
